# ansible/templates/benchtest.service.j2
[Unit]
Description=BenchTest InstantPayment API
After=network.target

[Service]
Type=exec
User=benchadmin
Group=benchadmin
WorkingDirectory=/opt/benchtest
ExecStart=/opt/benchtest/InstantPaymentBench.Api

Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:5000
Environment=DOTNET_ThreadPool_ForceMinWorkerThreads=1000
Environment=DOTNET_ThreadPool_ForceMinIOCompletionThreads=1000
Environment=DOTNET_gcServer=1
Environment=DOTNET_GCConcurrent=1

LimitNOFILE=65535
LimitNPROC=65535

Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15

SyslogIdentifier=benchtest

[Install]
WantedBy=multi-user.target
